# infrastructure/modules/ec2/systemd/nextjs-production.service
[Unit]
Description=Next.js Production Application
Documentation=https://nextjs.org/docs
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/srv/www/production/current
Environment=PATH=/home/ubuntu/.bun/bin:/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Environment=PORT=3002
ExecStart=/home/ubuntu/.bun/bin/bun start
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/srv/www/production
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nextjs-production

[Install]
WantedBy=multi-user.target