# infrastructure/modules/ec2/systemd/nextjs-staging.service
[Unit]
Description=Next.js Staging Application
Documentation=https://nextjs.org/docs
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/srv/www/staging/current
Environment=PATH=/home/ubuntu/.bun/bin:/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/home/ubuntu/.bun/bin/bun start
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/srv/www/staging /var/log/was-032
PrivateTmp=yes

# Logging
StandardOutput=append:/var/log/was-032/staging.log
StandardError=append:/var/log/was-032/staging-error.log
SyslogIdentifier=nextjs-staging

[Install]
WantedBy=multi-user.target